# docs/supervisor_config.ini
[program:openim_api]
directory=/var/www/openim
command=/var/www/openim/openim
user=root
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
environment=VIDEO2TEXT_ENV=online
stdout_logfile=/var/log/supervisor/openim.log
stdout_logfile_maxbytes=0
stdout_logfile_backups=0